# lxc-host/check-disk-state.cron.d
*/5 * * * * root  /usr/local/bin/check-disk-state.py >/tmp/check-disk-state.log 2>&1
