# lxc-host/check-disk-state.cron.d
*/5 * * * * ricardo  /usr/local/bin/check-disk-state.py >/dev/null 2>&1
