# lxc-host/duckdns.cron.d
*/15 * * * * ricardo  /usr/local/bin/duckdns.sh >/dev/null 2>&1
