# torrent/flexget.cron.d
@hourly  ricardo  /usr/local/bin/online && /usr/local/bin/flexget --cron execute
