# torrent/flexget.cron.d
0 *    * * *   ricardo  /usr/local/bin/flexget --loglevel warning -c config.yml execute --now
