# torrent/flexget.cron.d
@reboot  ricardo  /usr/local/bin/flexget daemon start -d
